# config/systemd/ben-led-release.service
[Unit]
Description=BEN LED pin release (free GPIO 12/13/16 from firmware boot indicator)
# Aucune dépendance — ce service doit tourner le plus tôt possible et libérer
# les pins LED RGB du firmware avant que les services BEN tentent leur PWM setup.
After=local-fs.target
# Avant les services BEN qui pilotent la LED RGB. Pas de Before=basic.target
# car ça créerait un cycle (multi-user.target wants ben-led-release ↔
# basic.target required by multi-user.target ↔ Before=basic.target).
Before=ben-network-check.service ben-tic-reader.service ben-lora-receiver.service ben-ble-provisioner.service

[Service]
Type=oneshot
# Petite séquence "welcome" : éteint le boot indicator (vert tenu par le
# firmware) → flash blanc 150ms (R+G+B haut) → off → release des pins pour
# que les services suivants puissent les claim via PWM.
ExecStart=/usr/bin/pinctrl set 13 op dl
ExecStart=/bin/sleep 0.3
ExecStart=/usr/bin/pinctrl set 12 op dh
ExecStart=/usr/bin/pinctrl set 13 op dh
ExecStart=/usr/bin/pinctrl set 16 op dh
ExecStart=/bin/sleep 0.15
ExecStart=/usr/bin/pinctrl set 12 op dl
ExecStart=/usr/bin/pinctrl set 13 op dl
ExecStart=/usr/bin/pinctrl set 16 op dl
ExecStart=/bin/sleep 0.15
ExecStart=/usr/bin/pinctrl set 12 ip
ExecStart=/usr/bin/pinctrl set 13 ip
ExecStart=/usr/bin/pinctrl set 16 ip
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
