# config/systemd/ben-lora-receiver.service
[Unit]
Description=BEN LoRa Receiver
After=basic.target
# Watchdog : le service se relance si le self-test radio échoue (radio figée OU Oops noyau détecté).
# ESCALADE REBOOT : si le restart ne suffit pas — Oops noyau (kernel tainted DIE) → IRQ RX morte,
# ré-init radio « GPIO busy » → restart-loop sans fin — seul un REBOOT clear l'IRQ. 3 restarts en
# 300 s → reboot. PROUVÉ NÉCESSAIRE par l'incident ben-0001 10/07 (Oops irq/160-lg, restart KO).
StartLimitIntervalSec=300
StartLimitBurst=3
StartLimitAction=reboot

[Service]
Type=simple
# WATCHDOG=1 envoyé par le process (self-test radio) → NotifyAccess=main pour l'accepter.
NotifyAccess=main
User=ben
# WorkingDirectory writable par `ben` requis : rpi-lgpio (drop-in replacement de
# RPi.GPIO pour kernel Trixie 6.12+) crée un FIFO de notification `.lgd-nfy*`
# en chemin RELATIF dans le CWD. Sans WorkingDirectory= le service hérite `/`
# → FileNotFoundError, le service crashe en boucle.
WorkingDirectory=/opt/ben/repo/src/pi/lora-receiver
ExecStart=/usr/bin/python3 /opt/ben/repo/src/pi/lora-receiver/main.py
Restart=always
RestartSec=10
# Le process pingue sd_notify(WATCHDOG=1) tant que le self-test radio (REG_VERSION==0x12)
# passe. S'il cesse (radio figée), systemd restart après WatchdogSec. > intervalle ping (30 s).
WatchdogSec=90

# PAS d'[Install]/WantedBy : démarré uniquement par ben-network-check.service
# quand le réseau est up. Voir l'explication dans ben-tic-reader.service.
