# config/systemd/wifi-watchdog.service
[Unit]
Description=WiFi connectivity watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi_watchdog.sh
